# services/sysvinit/haminegate
#!/bin/sh
### BEGIN INIT INFO
# Provides:          haminegate
# Required-Start:    $network $local_fs $remote_fs
# Required-Stop:     $network $local_fs $remote_fs
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: HAMineGate - Minecraft proxy with MOTD fallback
# Description:       Manages HAProxy and the MOTD fallback server
#                    for Minecraft server proxy/gatekeeping.
### END INIT INFO

HAMINEGATE_WRAPPER="${HAMINEGATE_WRAPPER:-/usr/local/bin/haminegate}"

case "${1:-}" in
    start|stop|restart|status|check)
        [ -x "$HAMINEGATE_WRAPPER" ] || {
            echo "HAMineGate wrapper not found at $HAMINEGATE_WRAPPER" >&2
            exit 1
        }
        exec "$HAMINEGATE_WRAPPER" "$@"
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status|check}"
        exit 1
        ;;
esac
